# deploy/systemd/yamls-waitress.service.example
[Unit]
Description=YAMLS portal waitress service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/yamls-portal
EnvironmentFile=/opt/yamls-portal/production.env
ExecStart=/opt/yamls-portal/.venv/bin/waitress-serve --host=127.0.0.1 --port=8080 wsgi:app
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
